# webadmins/config/supervisord.ini
[unix_http_server]
file=/app_shell/backup-platform/webadmins/supervisord.sock

[supervisord]
logfile=/app_shell/backup-platform/webadmins/log/supervisord.log 
logfile_maxbytes=50MB     
logfile_backups=10        
loglevel=info             
pidfile=/app_shell/backup-platform/webadmins/supervisord.pid
nodaemon=false              
minfds=1024               
minprocs=200             

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app_shell/backup-platform/webadmins/supervisord.sock



[program:uwsgi]
command=/app_shell/backup-platform/bin/uwsgi --ini /app_shell/backup-platform/webadmins/config/uwsgi.ini
stopsignal=QUIT
autostart=true
startretries=100
autorestart=true
stdout_logfile=/app_shell/backup-platform/webadmins/log/supervisord_uwsgi_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/app_shell/backup-platform/webadmins/log/supervisord_uwsgi_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10


[program:schedule]
command=/app_shell/backup-platform/bin/python -u /app_shell/backup-platform/webadmins/scheduler.py
stopsignal=QUIT
autostart=true
startretries=100
autorestart=true
stdout_logfile=/app_shell/backup-platform/webadmins/log/supervisord_schedule_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/app_shell/backup-platform/webadmins/log/supervisord_schedule_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

[program:celery]
command=/app_shell/backup-platform/bin/python /app_shell/backup-platform/webadmins/manage.py celery worker -n 4 --loglevel=info
stopsignal=QUIT
autostart=true
startretries=100
autorestart=true
stdout_logfile=/app_shell/backup-platform/webadmins/log/supervisord_celery_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/app_shell/backup-platform/webadmins/log/supervisord_celery_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
